# deploy/baker/centos-7-x64/etc/systemd/system/watson.service
[Unit]
Description=Datawire Baker Watson
ConditionPathExists=/etc/datawire/watson.conf
After=network.target

[Service]
ExecStart=/usr/bin/watson -c /etc/datawire/watson.conf
Restart=always

[Install]
WantedBy=default.target
